# scripts/launchd/com.andytsang.aitherapist.streamlit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.andytsang.aitherapist.streamlit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/andytsang/Documents/Project/心理咨詢agent/.venv/bin/python3</string>
        <string>-m</string>
        <string>streamlit</string>
        <string>run</string>
        <string>app.py</string>
        <string>--server.headless</string>
        <string>true</string>
        <string>--server.port</string>
        <string>8501</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/andytsang/Documents/Project/心理咨詢agent</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/tmp/streamlit.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/streamlit.log</string>
</dict>
</plist>
